# hypr-systemd-user/polkit-kde-agent.service
[Unit]
Description=Polkit KDE Agent

[Service]
ExecStart=/usr/lib/polkit-kde-authentication-agent-1
Restart=on-failure
